# supervisord.conf
; supervisor config file

[unix_http_server]
file=/var/run/supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=/opt/mda_workflow/logs/supervisord.log  ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10          ; (num of main logfile rotation backups;default 10)
loglevel=info               ; (log level;default info; others: debug,warn,trace)
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; use a unix:// URL  for a unix socket

[program:MDA-MAIN]
command=bash -c 'sleep 7200 && python3 /opt/mda_workflow/src/main.py'
autostart=true
autorestart=true
stderr_logfile=/opt/mda_workflow/logs/mda-main.err.log
stderr_logfile_maxbytes = 50MB
stderr_logfile_backups = 2
stdout_logfile=/opt/mda_workflow/logs/mda-main.out.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 2